# deploy/systemd/catmap-backend.service
[Unit]
Description=Campus Cat Association Map Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/catmap/backend
EnvironmentFile=/opt/catmap/backend/.env
ExecStart=/opt/catmap/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
